# samples/TPE-App/Python3/FtpToS3/ftp-to-s3.cron
*/15 * * * * /usr/local/bin/python -u /app/run.py
